# doorman.service
[Unit]
Description=Door unlocks the front door

[Service]
Type=simple
ExecStart=java -jar /usr/local/bin/doorman.jar /dev/ttyACM0